# examples/com.maniple.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.maniple</string>

    <!--
      Run maniple as an HTTP server on port 8766.
      Uses uvx to always run the latest published version.
      If you installed maniple globally (uv tool install maniple-mcp),
      you can replace the uvx line with the path to the maniple binary.
    -->
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/uvx</string>
      <string>maniple-mcp</string>
      <string>--http</string>
      <string>--port</string>
      <string>8766</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/you</string>

    <!--
      PATH must include the directory containing uvx.
      Adjust if uvx is installed elsewhere (e.g. /usr/local/bin).
    -->
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/you</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- iTerm2 Python API requires an Aqua login session -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    <key>ProcessType</key>
    <string>Interactive</string>

    <!-- Log output (create ~/.maniple/logs/ first) -->
    <key>StandardOutPath</key>
    <string>/Users/you/.maniple/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/you/.maniple/logs/stderr.log</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
  </dict>
</plist>
